# rootfs/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
priority=10
directory=/var/www/html
command=php-fpm
user=root
autostart=true
autorestart=true